# deploy/streamvault.service
[Unit]
Description=StreamVault - Next.js Web App
After=network.target

[Service]
Type=simple
User=streamvault
Group=streamvault
WorkingDirectory=/opt/streamvault
ExecStart=/home/streamvault/.bun/bin/bun .next/standalone/server.js
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Environment
EnvironmentFile=/opt/streamvault/.env.production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3000

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/streamvault/uploads /opt/streamvault/torrents /opt/streamvault/db /opt/streamvault/.next
PrivateTmp=true

# Resource limits (4-core, 8GB VPS)
LimitNOFILE=65536
MemoryMax=4G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=streamvault

[Install]
WantedBy=multi-user.target
